# config/supervisord.conf
[supervisord]
nodaemon=true

[program:telegram-bot]
command=python -u /app/bot.py
directory=/app
autostart=true
autorestart=true
stopsignal=TERM
stopasgroup=true
killasgroup=true

# Send Supervisor-managed logs to Docker logs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Pass token from container environment to the bot
environment=TELEGRAM_BOT_TOKEN=%(ENV_TELEGRAM_BOT_TOKEN)s

[program:obsidian-sync]
command=ob sync --continuous --path /vault
directory=/vault
autostart=%(ENV_OB_SYNC_AUTOSTART)s
autorestart=true
stopsignal=TERM
stopasgroup=true
killasgroup=true

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
